# vaisala-reset-rain.target
[Unit]
Description=Resets the Vaisala rain counter daily at 17:00
StopWhenUnneeded=yes
