# vaisala-reset-rain.target
[Unit]
Description=Resets the Vaisala rain counter daily at 12:00
StopWhenUnneeded=yes
